# pgcluu_collectd.timer
[Unit]
Description=Timer for pgcluu_collected - 1 Hour

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target


