# pgcluu_collectd.timer
[Unit]
Description=Timer for pgcluu_collected - 2 Hour

[Timer]
OnBootSec=10min
OnUnitActiveSec=2h

[Install]
WantedBy=timers.target


