# pgcluu.timer
[Unit]
Description=Run periodically pgcluu.service

[Timer]
OnBootSec=15min
OnUnitActiveSec=3h

[Install]
WantedBy=timers.target


